# background-backup.service
[Unit]
Description=Background backups

[Service]
ExecStart=/usr/local/bin/background-backup

[Install]
WantedBy=basic.target